# DigitalGarage-ui.service
# Digital Garage@2026
# Copy to /etc/systemd/system/digital_garage-ui.service
# Edit User, WorkingDirectory, and ExecStart paths to match your setup
[Unit]
Description=DigitalGarage UI
After=network.target

[Service]
Type=simple
# CHANGE THESE to match your user and install path:
User=YOURUSER
WorkingDirectory=/home/YOURUSER/digital_garage-ui
ExecStart=/home/YOURUSER/digital_garage-ui/venv/bin/uvicorn app:app --port 8000 --host 0.0.0.0
Restart=always
RestartSec=3
EnvironmentFile=-/home/YOURUSER/digital_garage-ui/.env

[Install]
WantedBy=multi-user.target
